# jobs/payment-jobs/schedules/run-routing-slip-tasks.dev.cron
# 5 minutes from 7 am to 6 pm PDT/PST
*/5 7-18 * * *
