# jobs/payment-jobs/schedules/run-direct-pay-automated-refunds.test.cron
# 11pm to 1 am PDT/PST
*/5 23,0-1 * * *
