# jobs/payment-jobs/schedules/run-create-cfs-accounts.test.cron
# every 10 minutes between 6 am to 9 pm PDT/PST 
*/10 6-21 * * *
