# jobs/payment-jobs/schedules/run-create-cfs-accounts.cron
# every 5 minutes between 6 am to 9 pm PDT/PST 
*/5 6-21 * * *