# jobs/payment-jobs/schedules/run-create-cfs-accounts.cron
# every 5 minutes between 8 am to 9 pm PDT/PST 
*/5 8-21 * * *
