# jobs/payment-jobs/schedules/run-routing-slip-tasks.test.cron
# every 5 minutes from 8 am to 9 pm PDT/PST 
*/5 8-21 * * *

