# jobs/payment-jobs/schedules/run-direct-pay-automated-refunds.dev.cron
# Every 5th minute PDT/PST from 7 am to 6 pm
*/5 7-18 * * *
